# packetirc.socket
[Unit]
Description=PacketIRC for LinBPQ

[Socket]
ListenStream=6001
Accept=yes

[Install]
WantedBy=sockets.target
